# mac/hotkey/com.gevensa.upload-hotkey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gevensa.upload-hotkey</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gevensa/.local/bin/upload-hotkey</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/gevensa/Library/Logs/upload-hotkey.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gevensa/Library/Logs/upload-hotkey.err.log</string>
</dict>
</plist>
